# roles/etcd/templates/etcd.service.j2
[Unit]
Description=etcd
Documentation=https://github.com/coreos

[Service]

[Install]
WantedBy=multi-user.target